# dockervolumes/haproxy1/etc/supervisor/conf.d/services.conf
[program:keepalived]
command=keepalived -n

[program:haproxy]
command=haproxy -db -f /etc/haproxy/haproxy.cfg
stdout_logfile=/var/log/haproxy.log
stdout_logfile_maxbytes=1GB
redirect_stderr=true
