# dockervolumes/keystone1/etc/supervisor/conf.d/services.conf
[program:keystone]
command=/usr/bin/keystone-all --log-dir /var/log/keystone/
user=keystone
