# dockervolumes/config1/etc/supervisor/conf.d/ifmap.conf
[program:ifmap]
command=/usr/bin/ifmap-server /etc/ifmap-server/ifmap.properties
stdout_logfile= /var/log/contrail/ifmap.log
stderr_logfile=/var/log/contrail/ifmap-err.log
user=contrail
