# dockervolumes/glance1/etc/supervisor/conf.d/services.conf
[program:glance-api]
command=/usr/bin/glance-api --log-dir /var/log/glance/
user=glance

[program:glance-registry]
command=/usr/bin/glance-registry --log-dir /var/log/glance/
user=glance
