# dockervolumes/galera1/etc/supervisor/conf.d/services.conf
[program:mysqld]
command=/usr/sbin/mysqld
autostart=false

;[program:puppet]
;command=

[program:bootstrap]
command=/bin/bash /usr/bin/mysqld_safe --wsrep_cluster_address=gcomm://
autostart=false
